# support-bot.service
[Unit]
Description=AI Support Bot Discord Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/support-bot
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/support-bot/output.log
StandardError=append:/var/log/support-bot/error.log

# Environment file with secrets
EnvironmentFile=/path/to/support-bot/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/support-bot/data

[Install]
WantedBy=multi-user.target
